# docker/wis2-gc.cron
0 0 * * * su -c "wis2-gc clean --days=$WIS2_GC_STORAGE_DATA_RETENTION_DAYS" > /proc/1/fd/1 2>/proc/1/fd/2
